# mosquitto-tailscale/rootfs/etc/s6-overlay/s6-rc.d/mosquitto/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Mosquitto & Tailscale
# Runs mosquitto
# ==============================================================================
declare -a options

options+=(-c /etc/mosquitto/mosquitto.conf)

bashio::log.info "Starting mosquitto MQTT broker..."
exec mosquitto "${options[@]}"
